# appwebsvr.conf
description "appwebsvr"
author  "Haitong Chen"

start on runlevel [234]
stop on runlevel [0156]

chdir /home/appadmin/product/sylz
exec python main.pyc
respawn
